# data/cookbooks/gecos_ws_mgmt/files/default/gecos-nm.conf
start on runlevel [2360]

script
    mv /etc/NetworkManager/system-connections/chef-conns/* /etc/NetworkManager/system-connections/
end script
