# data/cookbooks/gecos_ws_mgmt/templates/default/scripts_onshutdown.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:          script_launch
# Required-Start:    $local_fs $syslog
# Required-Stop:     $local_fs $syslog
# Default-Start:     
# Default-Stop:      6 0
# Short-Description: Starts the scripts launcher
# Description:       starts scripts and commands on shutdown
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DESC=scripts_launcher

set -e

. /lib/lsb/init-functions

case "$1" in
	stop)
		echo -n "Launching $DESC on shutdown: "
    <% @shutdown.each do |cmd| %>
      <%= cmd %>;
    <% end %>
		;;

	*)
		echo "Usage: $NAME {start}" >&2
		exit 1
		;;
esac

exit 0
