# data/cookbooks/gecos_ws_mgmt/templates/default/scripts_onstartup.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:          script_launch
# Required-Start:    $local_fs $syslog
# Required-Stop:     $local_fs $syslog
# Default-Start:     2 
# Default-Stop:        
# Short-Description: starts the scripts launcher
# Description:       starts scripts and commands on startup 
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DESC=scripts_launcher

set -e

. /lib/lsb/init-functions

case "$1" in
	start)
		echo -n "Launching $DESC on startup: " 
    <% @startup.each do |cmd| %>
      <%= cmd %>;
    <% end %>
		;;

	*)
		echo "Usage: $NAME {start}" >&2
		exit 1
		;;
esac

exit 0
